# freeradius/configs/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:freeradius]
directory=/etc/raddb/
command=bash -c 'radiusd -f ${DEBUG:+-X}'
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:certhandler]
directory=/
command=bash -c '/usr/local/bin/cert-renew.sh; sleep 86400'
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[group:all]
programs=freeradius,certhandler

# below allows supervisorctl to connect to the unix socket
[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
